# backend/config/supervisor-sataab-hotel.conf
# Supervisor Configuration for SATAAB Hotel Telegram Bot
# Location: /etc/supervisor/conf.d/sataab-hotel-worker.conf

[program:sataab-hotel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/sataab-hotel/backend/artisan queue:work redis --sleep=3 --tries=3 --timeout=90 --max-jobs=1000 --max-time=3600
autostart=true
autorestart=true
numprocs=4
redirect_stderr=true
stdout_logfile=/var/log/supervisor/sataab-hotel-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stopwaitsecs=3600
stopasgroup=true
killasgroup=true
user=www-data
environment=LARAVEL_ENV=production,LARAVEL_QUEUE=default

# Priority (lower number = higher priority)
priority=999

# Restart policy
startsecs=10
startretries=3

# Event listeners
eventlistener=unix:///tmp/supervisor.sock

---

# Scheduler Configuration
# Location: /etc/supervisor/conf.d/sataab-hotel-scheduler.conf

[program:sataab-hotel-scheduler]
process_name=%(program_name)s
command=php /var/www/sataab-hotel/backend/artisan schedule:work
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/sataab-hotel-scheduler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stopwaitsecs=3600
stopasgroup=true
killasgroup=true
user=www-data
environment=LARAVEL_ENV=production

# Priority (higher than workers)
priority=998

# Restart policy
startsecs=10
startretries=3

---

# High Priority Queue Worker
# Location: /etc/supervisor/conf.d/sataab-hotel-worker-high.conf

[program:sataab-hotel-worker-high]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/sataab-hotel/backend/artisan queue:work redis --queue=high --sleep=1 --tries=3 --timeout=60 --max-jobs=500 --max-time=1800
autostart=true
autorestart=true
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/sataab-hotel-worker-high.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stopwaitsecs=3600
stopasgroup=true
killasgroup=true
user=www-data
environment=LARAVEL_ENV=production,LARAVEL_QUEUE=high

# Priority (highest)
priority=997

# Restart policy
startsecs=10
startretries=3

---

# Low Priority Queue Worker
# Location: /etc/supervisor/conf.d/sataab-hotel-worker-low.conf

[program:sataab-hotel-worker-low]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/sataab-hotel/backend/artisan queue:work redis --queue=low --sleep=5 --tries=2 --timeout=120 --max-jobs=200 --max-time=7200
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/sataab-hotel-worker-low.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stopwaitsecs=3600
stopasgroup=true
killasgroup=true
user=www-data
environment=LARAVEL_ENV=production,LARAVEL_QUEUE=low

# Priority (lower than default)
priority=1000

# Restart policy
startsecs=10
startretries=3

---

# Event Listener for Failed Jobs
# Location: /etc/supervisor/conf.d/sataab-hotel-failed-listener.conf

[eventlistener:sataab-hotel-failed-listener]
buffer_size=100
events=PROCESS_STATE
result_handler=sataab_hotel_failed_handler

[eventlistener:sataab_hotel_failed_handler]
command=php /var/www/sataab-hotel/backend/artisan queue:failed
buffer_size=100
events=TICK_60

---

# Supervisor Group Configuration
# Add to /etc/supervisor/supervisord.conf

[group:sataab-hotel]
programs=sataab-hotel-worker,sataab-hotel-scheduler,sataab-hotel-worker-high,sataab-hotel-worker-low
priority=999

---

# Usage Commands:

# Reload supervisor configuration
# sudo supervisorctl reread
# sudo supervisorctl update

# Start all services
# sudo supervisorctl start sataab-hotel:*

# Stop all services
# sudo supervisorctl stop sataab-hotel:*

# Restart all services
# sudo supervisorctl restart sataab-hotel:*

# Check status
# sudo supervisorctl status

# View logs
# tail -f /var/log/supervisor/sataab-hotel-worker.log
# tail -f /var/log/supervisor/sataab-hotel-scheduler.log

# Restart specific worker
# sudo supervisorctl restart sataab-hotel-worker:0

# Start specific worker
# sudo supervisorctl start sataab-hotel-worker:0

# Stop specific worker
# sudo supervisorctl stop sataab-hotel-worker:0
