# backend/supervisor/hrms-worker.conf
# =============================================================================
# Supervisor configuration for HR&MS queue workers
# =============================================================================
# Install:
#   sudo cp supervisor/hrms-worker.conf /etc/supervisor/conf.d/hrms-worker.conf
#   sudo supervisorctl reread
#   sudo supervisorctl update
#   sudo supervisorctl start hrms-worker:*
#
# Monitor:
#   sudo supervisorctl status
# =============================================================================

[program:hrms-worker-default]
; Process label shown in supervisorctl
process_name=%(program_name)s_%(process_num)02d

; Adjust the path to your Laravel root
command=php /var/www/hrms/backend/artisan queue:work database \
    --queue=high,default,low \
    --sleep=3 \
    --tries=3 \
    --max-time=3600 \
    --timeout=90

; Number of parallel worker processes
numprocs=2

autostart=true
autorestart=true
startsecs=1
startretries=3

; Run as the web server user
user=www-data

; Redirect stderr to stdout so Supervisor captures everything
redirect_stderr=true
stdout_logfile=/var/log/supervisor/hrms-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5

; Stop signal — SIGTERM lets the worker finish its current job before stopping
stopsignal=TERM
stopwaitsecs=120
